# autoscaling/crontab
*/2 * * * * /scripts/scale-check.sh >> /logs/autoscale.log 2>&1
